# svxlink-relay.service
[Unit]
Description=SVXLink Reflector Relay
Documentation=man:socat(1)
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/svxlink-relay.sh --target TARGET_REFLECTOR_HOST --port 5300 --remote-port 5300 --log-file /var/log/svxlink-relay.log
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/run
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=svxlink-relay

# Environment
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target